# nfpm/rpm_files/postinst
if [ $1 -eq 1 ]; then
    /usr/bin/systemctl preset aingle.service >/dev/null 2>&1 ||:
    # print site info
    cat <<EOF
######################################
AINGLE installation finished
######################################
WARNING: The autopeering plugin will disclose your public IP address to possibly all nodes and entry points. Please disable the plugin prior to the first start if you do not want this to happen!

Please edit the config files (/var/lib/aingle/config.json and /var/lib/aingle/peering.json) to fit your needs.
To start AINGLE the first time, please run:

sudo systemctl enable aingle && sudo service aingle start

EOF

    mkdir -p /var/lib/aingle
    chmod 700 /var/lib/aingle
    chown aingle:aingle /var/lib/aingle

fi
